# deploy/systemd/user/tabura-dev-watch.service
[Unit]
Description=Watch Tabura source for dev reload
After=tabura-codex-app-server.service tabura-web.service
Wants=tabura-codex-app-server.service tabura-web.service

[Service]
Type=simple
WorkingDirectory=@@REPO_ROOT@@
ExecStart=@@REPO_ROOT@@/scripts/tabura-dev-watch.sh
Restart=always
RestartSec=1

[Install]
WantedBy=default.target
