# deploy/systemd/user/tabura-ptt.service
[Unit]
Description=Tabura Desktop Push-to-Talk Daemon
After=network.target tabura-stt.service
Wants=tabura-stt.service

[Service]
Type=simple
ExecStart=@@REPO_ROOT@@/tabura ptt-daemon
Restart=on-failure
RestartSec=5
SupplementaryGroups=input

[Install]
WantedBy=default.target
