# deploy/systemd/user/tabura-eou.service
[Unit]
Description=Tabura Semantic EOU Service
After=network.target

[Service]
Type=simple
WorkingDirectory=%h/code/assi/tabula
Environment=HOME=%h
Environment=TABURA_EOU_MODEL_REPO=livekit/turn-detector
Environment=TABURA_EOU_MODEL_DIR=%h/.local/share/tabura-eou/model
ExecStart=%h/.local/share/tabura-eou/venv/bin/uvicorn eou_server:app --app-dir %h/code/assi/tabula/scripts --host 127.0.0.1 --port 8425
Restart=on-failure
RestartSec=2

[Install]
WantedBy=default.target

