# deploy/systemd/user/tabura-piper-tts.service
[Unit]
Description=Tabura Piper TTS (fast, EN+DE)
After=network.target

[Service]
Type=simple
WorkingDirectory=%h/code/assi/tabula
Environment=HOME=%h
Environment=PIPER_MODEL_DIR=%h/.local/share/tabura-piper-tts/models
ExecStart=%h/.local/share/tabura-piper-tts/venv/bin/uvicorn piper_tts_server:app --app-dir %h/code/assi/tabula/scripts --host 127.0.0.1 --port 8424
Restart=on-failure
RestartSec=5
TimeoutStopSec=10

[Install]
WantedBy=default.target
